# embedded/embedded_programm/water_sensor.service
[Unit]
Description=Python Water Service
After=syslog.target
After=network.target

[Service]
Type=simple
PIDFile=/run/python_water_service.pid
ExecStart=
Restart=on-failure
OOMScoreAdjust=-100

[Install]
WantedBy=multi-user.target